# jetson/rc.local
echo -1 > /sys/module/usbcore/parameters/autosuspend
# /home/ubuntu/start_server_root.sh